# sync/com.roc.garminsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roc.garminsync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/matthewwalmsley/roc/sync/run.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/matthewwalmsley/roc/sync/sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/matthewwalmsley/roc/sync/sync.log</string>
</dict>
</plist>
